# config/crontab
# CRONJOBS
2 13 * * * source $HOME/.bashrc; venv/bin/python server/app/jobs/daily_bitcoin_price.py > /home/ubuntu/cron.log 2>&1







