# config/crontab
# CRONJOBS
SHELL=/bin/bash

15 13 * * * $HOME/venv/bin/python $HOME/server/app/jobs/daily_bitcoin_price.py > $HOME/cron.log 2>&1







