# config/crontab
# CRONJOBS
0 7 * * * venv/bin/python server/app/jobs/daily_bitcoin_price.py